# adapters/launchd/com.example.agentmem.coach.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.example.agentmem.coach</string>

  <key>ProgramArguments</key>
  <array>
    <!-- Confirm the node path matches `which node` on your machine. -->
    <string>/usr/local/bin/node</string>
    <string>/Users/YOUR_USER/Documents/code/agent-memory/bin/agentmem.mjs</string>
    <string>coach</string>
    <string>run</string>
  </array>

  <!-- Weekly on Sunday at 04:00 local time. Coaching is heavier than the
       nightly reflect pass (full knowledge base in the prompt) and the
       recommendations are aimed at a weekly review cadence anyway. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>0</integer>
    <key>Hour</key>
    <integer>4</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>EnvironmentVariables</key>
  <dict>
    <!-- API key is read from the macOS Keychain at runtime via the `security`
         CLI (service "agentmem", account "anthropic"). No plaintext on disk.
         See README for setup. Set AGENTMEM_API_KEY here only if you want to
         bypass the keychain entirely. -->
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/agentmem.coach.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/agentmem.coach.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
